# deploy/gaugeflow-social-telegram.service
[Unit]
Description=GaugeFlow Social OS — Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sam
WorkingDirectory=/home/sam/gaugeflow-social-os
EnvironmentFile=/home/sam/gaugeflow-social-os/.env
ExecStart=/home/sam/gaugeflow-social-os/.venv/bin/python main.py telegram
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/sam/gaugeflow-social-os/storage/logs/telegram.log
StandardError=append:/home/sam/gaugeflow-social-os/storage/logs/telegram.log

[Install]
WantedBy=multi-user.target
